# templates/haproxy-quadlet.container.j2
# HAProxy Blue-Green Deployment Quadlet Container Configuration
[Unit]
Description=HAProxy Blue-Green Load Balancer
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Container]
Image={{ haproxy.image }}
ContainerName={{ haproxy.container_name }}

# Port mappings
PublishPort={{ haproxy.http_port }}:80
PublishPort={{ haproxy.https_port }}:443
PublishPort={{ haproxy.stats_port }}:8404

# Volume mounts
{% for volume in podman_config.volumes %}
Volume={{ volume.source }}:{{ volume.target }}:{{ volume.type if volume.type is defined else 'bind' }}
{% endfor %}

# Network configuration
Network={{ podman_config.network_name }}

# Environment variables
Environment=HAPROXY_CONFIG=/usr/local/etc/haproxy/haproxy.cfg

# Security and resource limits
User=haproxy
Group=haproxy
ReadOnlyRootFilesystem=yes
NoNewPrivileges=yes
CapDrop=ALL
CapAdd=NET_BIND_SERVICE
CapAdd=CHOWN
CapAdd=SETGID
CapAdd=SETUID

# Logging
LogDriver=journald
LogOptions=tag="haproxy-blue-green"

# Health check
HealthCmd=haproxy -c -f /usr/local/etc/haproxy/haproxy.cfg
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=10s
HealthTimeout=10s

# Restart policy
Restart=on-failure
RestartSec=5s

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target 